# scripts/supervisord.conf
[program:flaskapp_server]
command=gunicorn app:app
numprocs=1
directory=/home/ubuntu/flaskapp/
user=ubuntu
